# apt-cache-proxy.service
[Unit]
Description=APT Cache Proxy Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/apt-cache-proxy
ExecStart=/usr/bin/python3 /opt/apt-cache-proxy/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
